# debian/cron.weekly
#!/bin/sh
# Used to update the documentation index for KDevelop

test -f "/usr/share/kdevelop/update_doc_index.sh" && "/usr/share/kdevelop/update_doc_index.sh" > /dev/null 2>&1
